# data_files/cosmicpi-ui.service
[Unit]
Description=Software for the CosmicPi UI

[Service]
ExecStart=/usr/bin/cosmicpi-ui
Restart=on-failure

[Install]
WantedBy=multi-user.target
