# install_files/CosmicPi-AP.service
[Unit]
Description=Start an AP for the CosmicPi

[Service]
ExecStart=/usr/bin/create_ap wlan0 eth0 CosmicPi

[Install]
WantedBy=multi-user.target
